# deploy/dealscanner-web.service
[Unit]
Description=DealScanner v3 web (next start, localhost only)
After=network.target

[Service]
User=dealscanner
WorkingDirectory=/opt/dealscanner-v2/web
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3

# Sandboxing: the web server needs no write access to the app tree at runtime.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
