# deploy/dealscanner-web.service
[Unit]
Description=DealScanner v3 web (next start, localhost only)
After=network.target

[Service]
User=dealscanner
WorkingDirectory=/opt/dealscanner-v2/web
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
