# deploy/dealscanner-api.service
[Unit]
Description=DealScanner v3 API (uvicorn, localhost only)
After=network.target

[Service]
User=dealscanner
WorkingDirectory=/opt/dealscanner-v2/api
EnvironmentFile=/opt/dealscanner-v2/.env
ExecStart=/opt/dealscanner-v2/engine/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8099
Restart=always
RestartSec=3

# Sandboxing: limit blast radius of any code-exec bug in the app.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/dealscanner-v2/data
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
